# crontab/crontab
# nightly updates

# brew
0 0 * * *   brew update
0 0 * * *   brew upgrade

# npm
0 0 * * *   npm update npm -g
0 0 * * *   npm update -g

# vim plugins
0 0 * * *   vim -c PluginUpdate
